# roles/common/templates/monthly-reboot.service.j2
[Unit]
Description=Monthly server reboot

[Service]
Type=oneshot
ExecStart=/sbin/reboot
