# roles/backup/templates/megakuponi-backup.service.j2
[Unit]
Description=Megakuponi restic backup

[Service]
Type=oneshot
ExecStart=/opt/megakuponi/scripts/backup.sh
